# deploy/linode-panel.service
[Unit]
Description=Linode Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=linode-panel
Group=linode-panel
EnvironmentFile=-/etc/linode-panel/linode-panel.env
ExecStart=/usr/local/bin/linode-panel
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/linode-panel

[Install]
WantedBy=multi-user.target
